# dist/etc/rc.local
#!/bin/sh -e

/home/pi/go.sh init
exit 0

